# debian/cron.d
0 3	* * *	root	test -x /usr/sbin/hlbackup && /usr/sbin/hlbackup
