# debian/cron.d
0 2	* * *	root	(test -x /usr/sbin/hlbackup && /usr/sbin/hlbackup); (test -x /usr/sbin/hlbackup-clean && /usr/sbin/hlbackup-clean /backup/*)
